# usr/lib/systemd/system/user@.service.d/privleapd.conf
# Copyright (C) 2025 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
# See the file COPYING for copying conditions.

[Service]
ExecStartPre=+/usr/bin/leapctl --create %i
ExecStopPost=+/usr/bin/leapctl --destroy %i
